# docker/workspace/crontab/hive
# Launch artisan CRONTAB every minutes
* * * * * hive /usr/bin/php /var/www/artisan schedule:run >> /dev/null 2>&1